# ext-sources/spamd.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">

<service_bundle type='manifest' name='spamd'>

<service
	name='network/spamd'
	type='service'
	version='1'>

	<create_default_instance enabled='false' />
	<single_instance />

	<dependency
	    name='fs-local'
	    grouping='require_all'
	    restart_on='none'
	    type='service'>
		<service_fmri value='svc:/system/filesystem/local' />
	</dependency>

	<dependency
	    name='network-service'
	    grouping='require_all'
	    restart_on='none'
	    type='service'>
		<service_fmri value='svc:/network/service' />
	</dependency>

	<dependent
		name='spamd_multi-user'
		grouping='optional_all'
		restart_on='none'>
			<service_fmri
			    value='svc:/milestone/multi-user' />
	</dependent>

	<exec_method
		type='method'
		name='start'
		exec='/usr/bin/spamd -i %{options/bind} -p %{options/port} -m %{options/max-children}'
		timeout_seconds='60' />

	<exec_method
		type='method'
		name='stop'
		exec=':kill'
		timeout_seconds='60' />

	<property_group name='startd' type='framework'>
		<propval name='duration' type='astring' value='child' />
	</property_group>

	<property_group name='options' type='application'>
		<propval name='bind' type='astring' value='127.0.0.1' />
		<propval name='port' type='integer' value='783' />
		<propval name='max-children' type='integer' value='5' />
	</property_group>

	<stability value='Unstable' />

	<template>
		<common_name>
			<loctext xml:lang='C'>
			Daemonized Version of SpamAssassin
			</loctext>
		</common_name>
		<documentation>
			<manpage title='spamd' section='1M'
			    manpath='/usr/share/man' />
		</documentation>
	</template>
</service>

</service_bundle>
